# files/etc/init.d/pos_service
#!/bin/sh /etc/rc.common

START=94

insmod /lib/modules/3.10.14/mt7620_gpio_int.ko
insmod /lib/modules/3.10.14/mt7620_gpio_key.ko
/usr/bin/pos_app &

